# catkin_ws/devel/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/multipos5/catkin_ws/src/ros/tools/rosbash/rostcsh"
