# catkin_ws/devel/etc/catkin/profile.d/15.rosbash.fish
# generated from rosbash/env-hooks/15.rosbash.fish.em

. "/home/multipos5/catkin_ws/src/ros/tools/rosbash/rosfish"
